# openvino_whisper/rootfs/etc/s6-overlay/s6-rc.d/ov_stt_engine/run
#!/usr/bin/with-contenv bashio
export MODEL_ID=$(bashio::config 'model')
export DEVICE=$(bashio::config 'device')
export LANGUAGE=$(bashio::config 'language')
exec python3 /app/server.py
