# openvino_whisper/rootfs/etc/s6-overlay/s6-rc.d/whisper_engine/run
#!/usr/bin/with-contenv bashio
export MODEL_ID=$(bashio::config 'model')
export DEVICE=$(bashio::config 'device')
export LANGUAGE=$(bashio::config 'language')

bashio::log.info "Starting OpenVINO Whisper Wyoming Server..."
exec python3 /app/server.py
